# .config/supervisor/local/supervisor-front.conf
[program:front]
command = npm run dev --prefix /srv/app/front